# claudecode-devos/systemd/devos-report.service
[Unit]
Description=ClaudeCode DevOS Report Generator

[Service]
Type=oneshot
User=kensan
WorkingDirectory=/opt/claudecode-devos
Environment=DEVOS_HOME=/opt/claudecode-devos
ExecStart=/usr/bin/python3 /opt/claudecode-devos/reports/report_generator.py
