# claudecode-devos/systemd/devos-dashboard.service
[Unit]
Description=ClaudeCode DevOS Flask Dashboard
After=network-online.target
Wants=network-online.target

[Service]
User=kensan
WorkingDirectory=/opt/claudecode-devos/web
Environment=DEVOS_HOME=/opt/claudecode-devos
Environment=DEVOS_PROJECTS_ROOT=/home/kensan/Projects
ExecStart=/usr/bin/gunicorn --workers 2 --bind 0.0.0.0:5050 --timeout 120 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
