# claudecode-devos/systemd/devos-evolution.timer
[Unit]
Description=Run ClaudeCode DevOS Evolution Engine hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
