# claudecode-devos/systemd/devos-report.timer
[Unit]
Description=Run ClaudeCode DevOS Report Generator hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
