# claudecode-devos/systemd/log-retention.timer
[Unit]
Description=Run ClaudeCode DevOS Log Retention daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
