# claudecode-devos/systemd/memory-guard.timer
[Unit]
Description=Run ClaudeCode DevOS Memory Guard every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Persistent=true
Unit=memory-guard.service

[Install]
WantedBy=timers.target
