# claudecode-devos/systemd/usage-manager.service
[Unit]
Description=ClaudeCode DevOS Usage Manager

[Service]
Type=oneshot
User=kensan
WorkingDirectory=/opt/claudecode-devos
Environment=DEVOS_HOME=/opt/claudecode-devos
ExecStart=/usr/bin/python3 /opt/claudecode-devos/ops/usage_manager.py reset
