# claudecode-devos/systemd/worker-heartbeat.service
[Unit]
Description=ClaudeCode DevOS Worker Heartbeat

[Service]
Type=oneshot
User=kensan
WorkingDirectory=/opt/claudecode-devos
Environment=DEVOS_HOME=/opt/claudecode-devos
ExecStart=/usr/bin/python3 /opt/claudecode-devos/cluster/workers/send_heartbeat.py
